# crontab
* * * * * root . /app/.env; echo running && /app/emailgrades.sh
# empty